# unit_files/entirem.timer
[Unit]
Description="Run entirem.service every two hours."

[Timer]
OnCalendar=00/2:10
Unit=entirem.service

[Install]
WantedBy=multi-user.target
