# unit_files/mtac.timer
[Unit]
Description="Run mtac.service every two hours."

[Timer]
OnCalendar=00/2:20
Unit=mtac.service

[Install]
WantedBy=multi-user.target
